# autostart/clickup.desktop
[Desktop Entry]
Exec=/usr/bin/clickup
Name=clickup
Type=Application
