# autostart/gajim.desktop
[Desktop Entry]
Exec=/usr/bin/gajim
Name=gajim
Comment=A fully-featured XMPP chat client
Type=Application
