# autostart/thunderbird.desktop
[Desktop Entry]
Exec=/usr/bin/thunderbird
Name=thunderbird
Type=Application
