# autostart/ticktick.desktop
[Desktop Entry]
Exec=/usr/bin/ticktick
Name=ticktick
Type=Application
