# autostart/TogglDesktop.desktop
[Desktop Entry]
Exec=/usr/bin/TogglDesktop
Name=TogglDesktop
Type=Application
